# deployment/personalcfo-celery.service
[Unit]
Description=PersonalCFO Celery Worker
After=network.target redis.service personalcfo.service
Wants=redis.service

[Service]
Type=exec
User=personalcfo
Group=personalcfo
WorkingDirectory=/opt/personalcfo
Environment=PATH=/opt/personalcfo/venv/bin
Environment=TMPDIR=/opt/personalcfo/tmp
ExecStartPre=/bin/mkdir -p /opt/personalcfo/tmp
ExecStart=/opt/personalcfo/venv/bin/celery -A app.core.celery_app.celery_app worker --loglevel=INFO -Q statements,priority
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=personalcfo-celery

# Security settings
NoNewPrivileges=true
ReadWritePaths=/opt/personalcfo

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
